# backup-system/systemd/ckis-backup.timer
[Unit]
Description=CKIS autonomous backup — reconciling safety net
Documentation=https://github.com/aedneth/ckis

[Timer]
# Frequent reconcile = the agent-agnostic real-time floor (no inotify dependency,
# survives any coding agent / OS / hardware). Interval + boot delay are templated
# from ckis-manifest.json (.schedule) by install.sh. Persistent catches missed
# runs after sleep/suspend; jitter avoids thundering herds.
OnBootSec=@BOOTDELAY@
OnUnitActiveSec=@INTERVAL@
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
